# bin/registry.service
[Unit]
Description=Worker Registry
Documentation=https://github.com/codex-team/registry
After=network.target

[Service]
Environment=NODE_ENV=production
Type=simple
User=registry
ExecStart=/home/registry/.nvm/versions/node/v11.2.0/bin/node /home/registry/registry/index.js
Restart=on-failure

[Install]
WantedBy=multi-user.target
